# .configs/etc/systemd/system/nodejs-watcher.service
[Unit]
Description=Som do Mato API Server Watcher

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart nodejs

[Install]
WantedBy=multi-user.target