# etc/systemd/system/sdm-watcher.path
[Path]
PathModified=/var/www/api.somdomato.com

[Install]
WantedBy=multi-user.target